# com.local.anchorman-wifi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Unique identifier for this daemon -->
    <key>Label</key>
    <string>com.local.anchorman-wifi</string>

    <!-- Path to the watcher script -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/usr/local/bin/anchorman-wifi</string>
    </array>

    <!-- Start automatically at boot -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart the daemon if it crashes or exits -->
    <key>KeepAlive</key>
    <true/>

    <!-- Log stdout and stderr -->
    <key>StandardOutPath</key>
    <string>/var/log/anchorman-wifi.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/anchorman-wifi.error.log</string>

    <!-- Runs as root so networksetup can toggle WiFi without prompts -->
    <!-- (LaunchDaemons in /Library/LaunchDaemons already run as root) -->

    <!-- Throttle restart attempts if the script keeps failing -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
